# packaging/linux/scripts/prerm
#!/bin/sh
set -e
case "$1" in
  remove|deconfigure)
    if [ -d /run/systemd/system ]; then
      systemctl disable --now dimmitd.service >/dev/null 2>&1 || true
    fi
    ;;
esac
exit 0
